# configs/automatic-backup/tuleap-deployment-backup.service
[Unit]
Description=Tuleap Deployment Backup

[Service]
ExecStart=/usr/bin/tuleap-deployment-backup.sh
